# lernstick-ltsp-user/DEBIAN/postinst
#! /bin/sh
# postinst script for erddcd
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> 'configure' <most-recently-configured-version>
#        * <old-postinst> 'abort-upgrade' <new version>
#        * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#          <failed-install-package> <version> 'removing'
#          <conflicting-package> <version>
# for details, see /usr/share/doc/packaging-manual/
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the 'postinst' is called with 'abort-upgrade',
#     'abort-remove' or 'abort-deconfigure'.

case "$1" in
    configure)

	systemctl enable  ltsp-update-config.service

	#no autologin
	if [ -z "$2" ] ; then
	#only on first install 
            if [ -e /var/lib/live/config/gdm3 ] ; then
                echo -n "Deactivating autologin for GDM ... "
                sed -i -e "s|^AutomaticLoginEnable=true\$|AutomaticLoginEnable=false|" \
                    -e "s|^TimedLoginEnable=true\$|TimedLoginEnable=false|" \
                    /etc/gdm3/daemon.conf
                echo "live-guest" > /var/lib/live/config/gdm3
                echo "done."
            fi
        fi

	#pam_mkhomedir

	#no updates
	dconf write /org/gnome/software/allow-updates false
	dconf write /org/gnome/software/download-updates false
	
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
